# scripts/crontab
0 3 * * * cd <path to repo>; ./scripts/my-builder.sh
