# scripts/crontab
0 3 * * * cd /home/iceflower/re-docker/ && ./scripts/my-builder.sh > /home/iceflower/cron.log 2>&1
